# systemd/did-i-notifications.timer
[Unit]
Description=Run did-i notifications

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target